# rootfs/usr/lib/systemd/user/steam-powerbuttond.service
[Unit]
Description=Steam Power Button Service

[Service]
ExecStart=/usr/bin/steam-powerbuttond
Restart=on-failure

[Install]
WantedBy=default.target
